# root/etc/s6-overlay/s6-rc.d/init-graphs1090-config/run
#!/usr/bin/with-contenv bash

set -euo pipefail

graphs1090_user=${GRAPHS1090_USER:-abc}
graphs1090_port=${GRAPHS1090_PORT:-8080}

if [[ "${graphs1090_user}" != "root" && "${graphs1090_user}" == "abc" ]]; then
    [[ -z ${PUID:-} ]] && export PUID=911
    [[ -z ${PGID:-} ]] && export PGID=911
fi

# Create required runtime directories
mkdir -p /config /run/graphs1090 /run/collectd /var/lib/collectd/rrd \
         /var/lib/graphs1090/scatter /usr/share/graphs1090/data-symlink \
         /usr/share/graphs1090/978-symlink \
         /tmp/nginx-client-body /tmp/nginx-proxy-temp \
         /tmp/nginx-fastcgi-temp /tmp/nginx-uwsgi-temp /tmp/nginx-scgi-temp 2>/dev/null || true

# Mount fontconfig cache to tmpfs to avoid unnecessary disk writes
if ! mount | grep -qs /var/cache/fontconfig 2>/dev/null; then
    mount -o rw,nosuid,nodev,relatime,size=32000k,mode=755 -t tmpfs tmpfs /var/cache/fontconfig 2>/dev/null || true
fi

# Install default config if not present in /config volume
if [[ ! -f /etc/default/graphs1090 ]]; then
    cp /usr/share/graphs1090/default-config /etc/default/graphs1090
fi

# Set up data symlink for collectd to read decoder JSON
DECODER_DIR="${GRAPHS1090_DECODER_DIR:-/run/readsb}"
ln -snf "${DECODER_DIR}" /usr/share/graphs1090/data-symlink/data

if [[ ! -d "${DECODER_DIR}" ]]; then
    printf '[init] WARNING: decoder directory %s not found — ensure a readsb/dump1090 container shares a volume there\n' "${DECODER_DIR}" >&2
elif [[ ! -f "${DECODER_DIR}/stats.json" ]]; then
    printf '[init] WARNING: %s/stats.json not found — ADS-B graphs will be incomplete until the decoder is running\n' "${DECODER_DIR}" >&2
fi

# Force collectd Hostname to "localhost" so RRD paths are predictable
# (critical for Balena and platforms with random container hostnames)
# NOTE: must be two separate sed calls — d immediately starts the next cycle,
# so combining '/^Hostname/d; 1s|^|...|' silently drops the Hostname line.
sed -i '/^Hostname/d' /etc/collectd/collectd.conf
sed -i '1i\Hostname "localhost"' /etc/collectd/collectd.conf

# Set up 978 UAT symlink if a UAT decoder directory is provided
UAT_DIR="${GRAPHS1090_UAT_DIR:-}"
if [[ -n "${UAT_DIR}" && -d "${UAT_DIR}" ]]; then
    ln -snf "${UAT_DIR}" /usr/share/graphs1090/978-symlink/data
    sed -i 's|^#\s*URL_978 .*|            URL_978 "file:///usr/share/graphs1090/978-symlink"|' /etc/collectd/collectd.conf
    printf '[init] UAT 978 MHz data directory: %s\n' "${UAT_DIR}"
fi

# Readback: restore RRD data from disk to /run/collectd (aligned with upstream readback.sh)
DBFOLDER=/var/lib/collectd/rrd
RUNFOLDER=/run/collectd

rm -f "${RUNFOLDER}/readback-complete"
rm -rf "${RUNFOLDER}/localhost"

# Safety net: if a previous run wrote data under a non-localhost hostname,
# migrate it so the current (hostname-locked) collectd can use it.
for _d in "${DBFOLDER}"/*/; do
    _name="$(basename "$_d")"
    if [[ "${_name}" != "localhost" && -d "${DBFOLDER}/${_name}" ]]; then
        printf '[init] migrating RRD data from hostname "%s" to "localhost"\n' "${_name}"
        mkdir -p "${DBFOLDER}/localhost"
        cp -a "${DBFOLDER}/${_name}/." "${DBFOLDER}/localhost/" 2>/dev/null || true
        rm -rf "${DBFOLDER}/${_name}"
    fi
done

readback_tar() {
    local tarfile="$1"
    if [[ ! -f "${tarfile}" ]]; then
        return 2
    fi
    # Skip files that are too small to be valid
    if (( $(stat -c %s "${tarfile}") < 50000 )); then
        printf '[init] WARNING: readback of %s skipped, file too small\n' "${tarfile}" >&2
        return 2
    fi
    if tar --overwrite --directory "${RUNFOLDER}" -xf "${tarfile}" 2>/dev/null; then
        printf '[init] readback of %s was successful\n' "${tarfile}"
        return 0
    fi
    printf '[init] ERROR: readback of %s failed\n' "${tarfile}" >&2
    rm -rf "${RUNFOLDER}/localhost"
    return 9
}

readback_folder() {
    local folder="$1"
    if [[ ! -d "${folder}" ]]; then
        return 2
    fi
    if cp -a "${folder}" "${RUNFOLDER}/localhost" 2>/dev/null; then
        printf '[init] readback of %s was successful\n' "${folder}"
        return 0
    fi
    printf '[init] ERROR: readback of %s failed\n' "${folder}" >&2
    return 9
}

readback_success() {
    mkdir -p "${RUNFOLDER}/localhost"
    # Delete empty RRD files (collectd chokes on them)
    find "${RUNFOLDER}/localhost" -size -50c -type f -delete 2>/dev/null || true
    touch "${RUNFOLDER}/readback-complete"
}

# Try primary tar, then folder, then weekly backups
if readback_tar "${DBFOLDER}/localhost.tar.gz"; then
    readback_success
elif readback_folder "${DBFOLDER}/localhost"; then
    readback_success
elif readback_tar "${DBFOLDER}/auto-backup-$(date +%Y-week_%V).tar.gz"; then
    readback_success
elif readback_tar "${DBFOLDER}/auto-backup-$(date +%Y-week_%V -d '1 week ago' 2>/dev/null).tar.gz" 2>/dev/null; then
    readback_success
else
    printf '[init] no valid existing RRD data found, starting fresh\n'
    readback_success
fi

# Create RRD stubs in the dump1090 directory.
# graphs1090's check() returns dump1090_dbfs-NaN.rrd for any missing file. Without
# a valid NaN.rrd, ALL graph generation fails. It must have only DS:value since
# collectd's dump1090.py plugin writes to it with a single data source.
# Additionally, dump1090_aircraft-recent.rrd uses DS names "total" and "positions"
# (not "value"), so check() falling back to NaN.rrd causes rrdtool errors for
# aircraft graphs. We create it as a separate stub with the correct DS names.
NANDIR="${RUNFOLDER}/localhost/dump1090-localhost"
mkdir -p "${NANDIR}"

NANRRD="${NANDIR}/dump1090_dbfs-NaN.rrd"
# Validate existing NaN.rrd has exactly 1 DS (value). A previous version created
# it with 3 DS, causing collectd "expected 3 data source readings (got 1)" errors.
# The readback tar preserves the bad file, so we must check and recreate.
_nan_ok=false
if [[ -f "${NANRRD}" ]]; then
    nan_ds_count=$(rrdtool info "${NANRRD}" 2>/dev/null | grep -c '\.type = ' || true)
    [[ "${nan_ds_count}" -eq 1 ]] && _nan_ok=true
fi
if ! "${_nan_ok}"; then
    rm -f "${NANRRD}"
    rrdtool create "${NANRRD}" --step 60 \
        DS:value:GAUGE:120:U:U \
        RRA:AVERAGE:0.8:1:3600 \
        RRA:MAX:0.8:1:3600 \
        RRA:MIN:0.8:1:3600
fi

ACFTRRD="${NANDIR}/dump1090_aircraft-recent.rrd"
if [[ ! -f "${ACFTRRD}" ]]; then
    rrdtool create "${ACFTRRD}" --step 60 \
        DS:total:GAUGE:120:U:U \
        DS:positions:GAUGE:120:U:U \
        RRA:AVERAGE:0.8:1:3600 \
        RRA:MAX:0.8:1:3600 \
        RRA:MIN:0.8:1:3600
fi

# Create DERIVE-type stubs for tracks RRDs. Without these, the graph
# check() function falls back to dump1090_dbfs-NaN.rrd (GAUGE=1), which
# the tracks graph misinterprets as ~3600 tracks/hour (1/s × 3600).
for _tracks_inst in all single_message; do
    _tracks_rrd="${NANDIR}/dump1090_tracks-${_tracks_inst}.rrd"
    if [[ ! -f "${_tracks_rrd}" ]]; then
        rrdtool create "${_tracks_rrd}" --step 60 \
            DS:value:DERIVE:120:0:16000 \
            RRA:AVERAGE:0.8:1:3600 \
            RRA:MAX:0.8:1:3600 \
            RRA:MIN:0.8:1:3600
    fi
done
printf '[init] RRD stubs ready\n'

# Diagnostic: show what directories exist after readback
printf '[init] collectd data directories in %s/localhost/:\n' "${RUNFOLDER}"
ls -1 "${RUNFOLDER}/localhost/" 2>/dev/null || printf '  (none)\n'

# Count RRD files for visibility
_rrd_total=$(find "${RUNFOLDER}/localhost" -name '*.rrd' 2>/dev/null | wc -l)
printf '[init] total RRD files after readback: %d\n' "${_rrd_total}"

# Symlink /var/lib/collectd/rrd/localhost → /run/collectd/localhost
# Upstream graphs1090.sh sources /etc/default/graphs1090 twice; the second source
# resets DB back to /var/lib/collectd/rrd after the autodetect. The symlink makes
# both paths equivalent so graphs work regardless of which DB value is active.
if [[ -d "${DBFOLDER}/localhost" && ! -L "${DBFOLDER}/localhost" ]]; then
    rm -rf "${DBFOLDER}/localhost"
fi
ln -snf "${RUNFOLDER}/localhost" "${DBFOLDER}/localhost"
printf '[init] symlink: %s/localhost → %s/localhost\n' "${DBFOLDER}" "${RUNFOLDER}"

# Belt-and-suspenders: also fix DB in config
sed -i 's|^DB=.*|DB=/run/collectd|' /etc/default/graphs1090

# Render nginx config from template
mkdir -p /run/nginx/http.d 2>/dev/null || true
sed -e "s|\${GRAPHS1090_PORT}|${graphs1090_port}|" \
    /etc/nginx/http.d/default.conf.template \
    > /run/nginx/http.d/default.conf

# Run boot.sh to configure HTML (nographs = skip graph generation, done by svc-graphs1090)
source /etc/default/graphs1090
if ! /usr/share/graphs1090/boot.sh nographs; then
    printf '[init] ERROR: boot.sh exited with errors\n' >&2
fi

# Copy placeholder PNGs for graphs that would otherwise 404 (e.g. misc when
# no gain-reporting hardware exists, disk_io when block devices aren't visible).
# The PNG is baked into the image at build time so this is a simple cp — no
# rrdtool dependency at runtime.
_placeholder_src="/usr/share/graphs1090/no-data.png"
_placeholder_count=0

if [[ ! -f "${_placeholder_src}" ]]; then
    printf '[init] WARNING: placeholder source %s not found — generating fallback\n' "${_placeholder_src}" >&2
    # Generate a minimal white PNG placeholder using Python (no rrdtool dependency)
    python3 -c "import struct,zlib;W=619;H=324;raw=b''.join(b'\x00'+bytes([255,255,255])*W for _ in range(H));c=zlib.compress(raw);chunk=lambda t,d:struct.pack('>I',len(d))+t+d+struct.pack('>I',zlib.crc32(t+d)&0xffffffff);f=open('${_placeholder_src}','wb');f.write(b'\x89PNG\r\n\x1a\n'+chunk(b'IHDR',struct.pack('>IIBBBBB',W,H,8,2,0,0,0))+chunk(b'IDAT',c)+chunk(b'IEND',b''));f.close()" 2>/dev/null || true
fi

# Ensure the placeholder is also in the HTML dir for nginx error_page fallback
if [[ -f "${_placeholder_src}" && ! -f /usr/share/graphs1090/html/no-data.png ]]; then
    cp "${_placeholder_src}" /usr/share/graphs1090/html/no-data.png
    printf '[init] copied placeholder to HTML dir for nginx fallback\n'
fi

if [[ -f "${_placeholder_src}" ]]; then
    printf '[init] placeholder source: %s (%s bytes)\n' "${_placeholder_src}" "$(stat -c %s "${_placeholder_src}" 2>/dev/null || echo '?')"

    # Graph basenames that may 404 when their data source is unavailable:
    #   dump1090-localhost-misc      — requires gain-reporting hardware
    #   system-localhost-disk_io_iops   — requires visible block devices
    #   system-localhost-disk_io_octets — requires visible block devices
    _placeholder_graphs=(
        "dump1090-localhost-misc"
        "system-localhost-disk_io_iops"
        "system-localhost-disk_io_octets"
    )

    for _graph in "${_placeholder_graphs[@]}"; do
        for p in 2h 8h 24h 48h 7d 14d 30d 90d 180d 365d 730d 1095d 1825d 3650d; do
            _target="/run/graphs1090/${_graph}-${p}.png"
            if [[ ! -f "${_target}" ]]; then
                if cp "${_placeholder_src}" "${_target}"; then
                    _placeholder_count=$(( _placeholder_count + 1 ))
                else
                    printf '[init] ERROR: failed to create placeholder %s\n' "${_target}" >&2
                fi
            fi
        done
    done
else
    printf '[init] ERROR: no placeholder source available — some graphs will 404 until data is collected\n' >&2
fi
printf '[init] placeholder PNGs: %d created\n' "${_placeholder_count}"

# Set ownership
if [[ -z ${LSIO_NON_ROOT_USER:-} ]]; then
    lsiown -R "${graphs1090_user}:${graphs1090_user}" \
        /config /run/nginx /run/graphs1090 /run/collectd \
        /var/lib/collectd /var/lib/graphs1090 \
        /usr/share/graphs1090 2>/dev/null || true
fi

printf '[init] graphs1090 initialization complete\n'
