# root/etc/s6-overlay/s6-rc.d/svc-graphs1090/run
#!/usr/bin/with-contenv bash

set -euo pipefail

# Wait for collectd to create at least one system RRD file before generating graphs.
# collectd uses a 60-second interval, so the first write occurs ~60-70s after start.
# If readback restored data from a previous run, files already exist and we proceed
# immediately.
RUNFOLDER=/run/collectd/localhost
WAIT_MAX=120
WAIT=0
while [[ ${WAIT} -lt ${WAIT_MAX} ]]; do
    # Check for aggregation-cpu-average (created by collectd's aggregation plugin)
    # or any system_stats directory (created by the Python system_stats plugin)
    if [[ -d "${RUNFOLDER}/aggregation-cpu-average" ]] || [[ -d "${RUNFOLDER}/system_stats" ]]; then
        break
    fi
    if [[ ${WAIT} -eq 0 ]]; then
        printf '[graphs1090] waiting for collectd to write initial RRD data...\n'
        # Diagnostic: show what exists so far
        if [[ -d "${RUNFOLDER}" ]]; then
            printf '[graphs1090] current contents of %s:\n' "${RUNFOLDER}"
            ls -1 "${RUNFOLDER}" 2>/dev/null || printf '  (empty)\n'
        else
            printf '[graphs1090] WARNING: %s does not exist yet\n' "${RUNFOLDER}"
        fi
    fi
    sleep 5
    WAIT=$((WAIT + 5))
done

if [[ ${WAIT} -ge ${WAIT_MAX} ]]; then
    printf '[graphs1090] WARNING: timed out waiting for RRD data after %ds, starting anyway\n' "${WAIT_MAX}"
    printf '[graphs1090] contents of %s:\n' "${RUNFOLDER}"
    ls -1 "${RUNFOLDER}" 2>/dev/null || printf '  (directory missing)\n'
else
    printf '[graphs1090] collectd RRD data detected after %ds\n' "${WAIT}"
fi

printf '[graphs1090] starting graph generation service\n'

# === Diagnostic: verify all prerequisites for graph generation ===
# This section runs ONCE at service start to expose hidden errors.
# The upstream scripts suppress all output with &>/dev/null.
printf '[graphs1090] === DIAGNOSTIC START ===\n'

# 1. Check GNU date (coreutils) - needed for END_TIME=$(date -d -1min ...)
printf '[graphs1090] date binary: %s\n' "$(command -v date 2>&1)"
if date -d '-1 minute' '+%H:%M' >/dev/null 2>&1; then
    printf '[graphs1090] GNU date -d test: OK (%s)\n' "$(date -d '-1 minute' '+%H:%M')"
else
    printf '[graphs1090] GNU date -d test: FAILED\n'
    printf '[graphs1090]   date --version: %s\n' "$(date --version 2>&1 | head -1)"
fi

# 2. Check renice capability (needs SYS_NICE, may not be granted)
if renice -n 20 -p $$ >/dev/null 2>&1; then
    printf '[graphs1090] renice: OK\n'
else
    printf '[graphs1090] renice: FAILED (no SYS_NICE cap - this is expected)\n'
fi

# 3. Check /run/graphs1090 is writable
if touch /run/graphs1090/_diag_test 2>/dev/null; then
    rm -f /run/graphs1090/_diag_test
    printf '[graphs1090] /run/graphs1090 writable: YES\n'
else
    printf '[graphs1090] /run/graphs1090 writable: NO\n'
    ls -la /run/ 2>&1 | grep -E 'graphs1090|^total' || true
fi

# 4. Check rrdtool can create a PNG from real data
source /etc/default/graphs1090
DB=/var/lib/collectd/rrd
if [[ -d /run/collectd/localhost ]]; then DB=/run/collectd; fi
printf '[graphs1090] DB=%s\n' "${DB}"

TESTRRD="${DB}/localhost/aggregation-cpu-average/cpu-idle.rrd"
if [[ -f "${TESTRRD}" ]]; then
    printf '[graphs1090] test RRD: %s (%s bytes)\n' "${TESTRRD}" "$(stat -c %s "${TESTRRD}" 2>/dev/null)"
    rrdtool graph /run/graphs1090/_diag.png \
        --imgformat PNG --end now --start end-1h \
        --width 100 --height 50 --title "diag-test" \
        "DEF:v=${TESTRRD}:value:AVERAGE" "LINE1:v#FF0000:test" 2>&1 || true
    if [[ -f /run/graphs1090/_diag.png ]]; then
        printf '[graphs1090] rrdtool basic graph: OK (%s bytes)\n' "$(stat -c %s /run/graphs1090/_diag.png)"
        rm -f /run/graphs1090/_diag.png
    else
        printf '[graphs1090] rrdtool basic graph: FAILED (no PNG created)\n'
    fi
else
    printf '[graphs1090] test RRD missing: %s\n' "${TESTRRD}"
    printf '[graphs1090] aggregation-cpu-average contents:\n'
    ls -la "${DB}/localhost/aggregation-cpu-average/" 2>&1 | head -5 || true
fi

# 5. Run graphs1090.sh 24h with stderr+stdout visible to expose rrdtool errors
printf '[graphs1090] running: /usr/share/graphs1090/graphs1090.sh 24h 0\n'
/usr/share/graphs1090/graphs1090.sh 24h 0 2>&1 | head -80 || true
printf '[graphs1090] graphs1090.sh exit: %d\n' "${PIPESTATUS[0]}"

# 6. Check resulting PNGs
png_count=$(find /run/graphs1090/ -name '*.png' 2>/dev/null | wc -l)
printf '[graphs1090] PNG files in /run/graphs1090/: %d\n' "${png_count}"
ls -la /run/graphs1090/*.png 2>&1 | head -10 || printf '  (none)\n'
printf '[graphs1090] === DIAGNOSTIC END ===\n'

exec /usr/share/graphs1090/service-graphs1090.sh
