# root/etc/s6-overlay/s6-rc.d/svc-graphs1090/run
#!/usr/bin/with-contenv bash

set -euo pipefail

# Wait for collectd to create at least one system RRD before generating graphs.
# If readback restored data, files already exist and we skip the wait.
RUNFOLDER=/run/collectd/localhost
WAIT_MAX=120
WAIT=0
while [[ ${WAIT} -lt ${WAIT_MAX} ]]; do
    if [[ -d "${RUNFOLDER}/aggregation-cpu-average" || -d "${RUNFOLDER}/system_stats" ]]; then
        break
    fi
    if [[ ${WAIT} -eq 0 ]]; then
        printf '[graphs1090] waiting for collectd RRD data...\n'
    fi
    sleep 5
    WAIT=$((WAIT + 5))
done

if [[ ${WAIT} -ge ${WAIT_MAX} ]]; then
    printf '[graphs1090] WARNING: timed out after %ds, starting anyway\n' "${WAIT_MAX}"
else
    printf '[graphs1090] collectd RRD data detected after %ds\n' "${WAIT}"
fi

printf '[graphs1090] starting graph generation service\n'

exec /usr/share/graphs1090/service-graphs1090.sh
