# root/etc/s6-overlay/s6-rc.d/svc-graphs1090/run
#!/usr/bin/with-contenv bash

set -euo pipefail

printf '[graphs1090] starting graph generation service\n'
exec /usr/share/graphs1090/service-graphs1090.sh
