# root/etc/s6-overlay/s6-rc.d/svc-graphs1090/run
#!/usr/bin/with-contenv bash

set -euo pipefail

# shellcheck disable=SC1091
source /usr/local/bin/log-functions.sh
SVC_NAME="graphs1090"

# Wait for collectd to create at least one system RRD before generating graphs.
# If readback restored data, files already exist and we skip the wait.
RUNFOLDER=/run/collectd/localhost
WAIT_MAX=120
WAIT=0
while [[ ${WAIT} -lt ${WAIT_MAX} ]]; do
    if [[ -d "${RUNFOLDER}/aggregation-cpu-average" || -d "${RUNFOLDER}/system_stats" ]]; then
        break
    fi
    if [[ ${WAIT} -eq 0 ]]; then
        log_info "waiting for collectd RRD data..."
    fi
    sleep 5
    WAIT=$((WAIT + 5))
    # Periodic status while waiting
    if (( WAIT % 15 == 0 )); then
        _dirs=$(ls -1 "${RUNFOLDER}/" 2>/dev/null | tr '\n' ' ' || echo 'none')
        log_info "still waiting (${WAIT}s)... directories so far: ${_dirs:-none}"
    fi
done

if [[ ${WAIT} -ge ${WAIT_MAX} ]]; then
    log_warn "timed out after ${WAIT_MAX}s waiting for collectd data, starting anyway"
else
    log_info "collectd RRD data detected after ${WAIT}s"
fi

# Log RRD directories and graph output state
_rrd_dirs=$(ls -1 "${RUNFOLDER}/" 2>/dev/null | tr '\n' ' ' || echo 'none')
log_info "RRD data directories: ${_rrd_dirs:-none}"

# Check for disk RRD data (helps diagnose disk_io graph issues)
_disk_rrds=$(find "${RUNFOLDER}" -path '*/disk-*/disk_*' -name '*.rrd' 2>/dev/null | head -5)
if [[ -n "${_disk_rrds}" ]]; then
    log_info "disk RRD data found:"
    echo "${_disk_rrds}" | while read -r f; do printf '  %s\n' "$f"; done
else
    log_warn "no disk RRD data found — disk_io graphs will show placeholder"
fi

_graph_count=$(find /run/graphs1090/ -name '*.png' 2>/dev/null | wc -l)
log_info "existing graph PNGs: ${_graph_count}"
log_info "starting graph generation service"

exec /usr/share/graphs1090/service-graphs1090.sh
