# root/etc/s6-overlay/s6-rc.d/svc-nginx/run
#!/usr/bin/with-contenv bash

set -euo pipefail

# shellcheck disable=SC1091
source /usr/local/bin/log-functions.sh
SVC_NAME="nginx"

graphs1090_port=${GRAPHS1090_PORT:-8080}

log_info "starting graphs1090 web UI on port ${graphs1090_port}"
log_info "web interface: http://<host-ip>:${graphs1090_port}/graphs1090/"

exec nginx -g 'daemon off;'
