# root/etc/s6-overlay/s6-rc.d/svc-collectd/finish
#!/usr/bin/with-contenv bash

# On collectd shutdown, write RRD data from /run/collectd back to persistent volume
# Aligned with upstream writeback.sh logic

# shellcheck disable=SC1091
source /usr/local/bin/log-functions.sh
SVC_NAME="collectd-finish"

TARGET=/var/lib/collectd/rrd
RUNFOLDER=/run/collectd

if [[ ! -f "${RUNFOLDER}/readback-complete" ]]; then
    log_warn "readback did not complete, skipping writeback"
    exit 0
fi

if [[ ! -d "${RUNFOLDER}/localhost" ]]; then
    log_warn "no data in ${RUNFOLDER}/localhost to write back"
    exit 0
fi

log_info "writing RRD data from ${RUNFOLDER} to disk..."
_rrd_count=$(find "${RUNFOLDER}/localhost" -name '*.rrd' 2>/dev/null | wc -l)
log_info "RRD files to write back: ${_rrd_count}"
mkdir -p "${TARGET}"

TMPF="${TARGET}/localhost.tar.gz.tmp"
rm -f "${TMPF}"

if tar --directory "${RUNFOLDER}" -c localhost | gzip -1 -c > "${TMPF}" && [[ -f "${TMPF}" ]]; then
    # Rotate existing tar.gz to weekly backup (only if large enough to be valid)
    if [[ -f "${TARGET}/localhost.tar.gz" ]] && (( $(stat -c %s "${TARGET}/localhost.tar.gz") > 150000 )); then
        BACKUP="${TARGET}/auto-backup-$(date +%Y-week_%V).tar.gz"
        # Only overwrite backup if it's not much larger than current
        if ! [[ -f "${BACKUP}" ]] || (( $(stat -c %s "${BACKUP}") < $(stat -c %s "${TARGET}/localhost.tar.gz") + 512 * 1024 )); then
            mv -f "${TARGET}/localhost.tar.gz" "${BACKUP}" 2>/dev/null || true
        fi
        # Remove backups older than 60 days
        find "${TARGET}" -name 'auto-backup-*.tar.gz' -mtime +60 -delete 2>/dev/null || true
    fi

    if ! sync "${TMPF}" 2>/dev/null; then
        log_warn "sync failed"
    fi
    mv -f "${TMPF}" "${TARGET}/localhost.tar.gz"
    log_info "writeback complete: $(du -sh "${TARGET}/localhost.tar.gz" 2>/dev/null || echo 'unknown size')"

    # Remove localhost folder (tar.gz is the canonical source)
    if [[ -d "${TARGET}/localhost" ]]; then
        rm -rf "${TARGET}/localhost"
    fi
else
    log_error "writeback failed"
    rm -f "${TMPF}"
fi
