# root/etc/s6-overlay/s6-rc.d/svc-collectd/run
#!/usr/bin/with-contenv bash

set -euo pipefail

printf '[collectd] starting collectd data collector\n'
exec collectd -f -C /etc/collectd/collectd.conf
