# root/etc/s6-overlay/s6-rc.d/svc-collectd/run
#!/usr/bin/with-contenv bash

set -euo pipefail

# shellcheck disable=SC1091
source /usr/local/bin/log-functions.sh
SVC_NAME="collectd"

log_info "starting collectd data collector"

# Log key collectd configuration for debugging
if [[ -f /etc/collectd/collectd.conf ]]; then
    _hostname=$(grep -m1 '^Hostname' /etc/collectd/collectd.conf | awk '{print $2}' || echo '?')
    _datadir=$(grep -m1 'DataDir' /etc/collectd/collectd.conf | awk '{print $2}' || echo '?')
    log_info "config: Hostname=${_hostname} DataDir=${_datadir}"
fi

# Log available block devices (helps diagnose disk_io 404s)
_blockdevs=$(lsblk -dnpo NAME 2>/dev/null | tr '\n' ' ' || echo 'none detected')
log_info "block devices visible: ${_blockdevs:-none detected}"

# Log decoder symlink status
_decoder_link=$(readlink -f /usr/share/graphs1090/data-symlink/data 2>/dev/null || echo 'missing')
log_info "decoder data symlink: ${_decoder_link}"
if [[ -f "${_decoder_link}/stats.json" ]]; then
    log_info "stats.json: found ($(stat -c %s "${_decoder_link}/stats.json" 2>/dev/null || echo '?') bytes)"
else
    log_warn "stats.json not found — ADS-B data will not be collected"
fi

exec collectd -f -C /etc/collectd/collectd.conf
